# templates/upstart.uwsgi.conf.j2
# Generated by Ansible lciolecki.uwsgi role
# uWSGI upstart script

description "uWSGI application script"
author "Łukasz Ciołecki"

start on filesystem and net-device-up IFACE!=lo
stop on runlevel [016]
respawn

setuid {{ uwsgi_user }}
setgid {{ uwsgi_group }}

kill timeout {{ uwsgi_upstart_timeout_kill }}

exec {{ uwsgi_source }} --ini {{ uwsgi_config_dir }}/uwsgi.ini
